# etc/supervisor.d/rpcstatus.conf
[program:rpcstatus]
command=/usr/local/bin/statusservice.py
user=cosmovisor
autostart=true
autorestart=unexpected
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0